# files/rc.local
#!/bin/bash

if [ -e /opt/ran_customization ]; then
    exit
else
    /root/setup.sh &> /var/log/photon-customization.log
fi
